# vps-cookie-refresher/noir-cookie-refresh.timer
[Unit]
Description=Noir - dispara el refresco de cookie cada 30 minutos

[Timer]
# Cada 30 min: mantiene la sesion bien fresca.
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
