# vps-cookie-refresher/noir-cookie-refresh.service
[Unit]
Description=Noir - refresco de cookie de YouTube Music
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/app/vps-cookie-refresher
User=root
Environment=BASE_URL=https://84-247-174-216.sslip.io
# JWT_SECRET (el mismo del backend) va en refresh.env con permisos 600.
EnvironmentFile=/root/app/vps-cookie-refresher/refresh.env
ExecStart=/usr/bin/node refresh.mjs
